# com.transcription-expert.health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.transcription-expert.health</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/macmini/apps/transcription-expert/venv/bin/python3</string>
        <string>/Users/macmini/apps/transcription-expert/health_check.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/macmini/apps/transcription-expert</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/macmini/apps/transcription-expert/health.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/macmini/apps/transcription-expert/health.log</string>
</dict>
</plist>
