# com.transcription-expert.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.transcription-expert</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/macmini/apps/transcription-expert/venv/bin/python3</string>
        <string>/Users/macmini/apps/transcription-expert/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/macmini/apps/transcription-expert</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/macmini/apps/transcription-expert/transcription.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/macmini/apps/transcription-expert/transcription.log</string>
</dict>
</plist>
